# packer/files/puma.service
[Unit]
Description=puma service
[Service]
WorkingDirectory=/opt/reddit
ExecStart=/usr/local/bin/puma
[Install]
WantedBy=multi-user.target
